# service/sprt-dev-restart.service
[Unit]
Description=Restart sprt-dev.service on file change

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart sprt-dev.service